# local.restaurant-hiring-monitor-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.restaurant-hiring-monitor-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_PATH__/watchdog.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_PATH__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Daily 10:00 — runs an hour after the scrape window so a same-day
             scrape success will already have updated the heartbeat. -->
        <key>Hour</key><integer>10</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <!-- Fallback every 6 hours so stale-data alerts are not missed if the
         calendar trigger is skipped while the Mac is asleep. -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <key>StandardOutPath</key>
    <string>__PROJECT_PATH__/logs/watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_PATH__/logs/watchdog.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
